# upgrade-notify-0330.timer
[Unit]
Description=Notify users at 03:30 before security upgrade
[Timer]
OnCalendar=Sun *-*-* 03:30:00
Persistent=true
[Install]
WantedBy=timers.target
